# playbook/common/roles/modules/logstash/files/common/etc/profile.d/logstash.sh
# LOGSTASH Profile

LOGSTASH_HOME="/opt/logstash"
PATH="${PATH}:${LOGSTASH_HOME}/bin"

export LOGSTASH_HOME PATH
